# core/cluster/below/com.buildrunner.below-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.buildrunner.below-worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/byronhudson/Projects/BuildRunner3/.venv/bin/python</string>
        <string>-m</string>
        <string>core.cluster.below.research_worker</string>
        <string>--queue-dir</string>
        <string>.buildrunner/research-queue</string>
        <string>--poll-seconds</string>
        <string>5</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/byronhudson/Projects/BuildRunner3</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/byronhudson/Library/Logs/br3-below-worker.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/byronhudson/Library/Logs/br3-below-worker.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>BR3_BELOW_WORKER_LOG_LEVEL</key>
        <string>INFO</string>
    </dict>
</dict>
</plist>
